# target_overlay/etc/profile.d/term.sh
# Makes mc run via ssh session
export TERM=vt100

